# .config/systemd/user/atuin.service
[Unit]
Description=Atuin Shell History Daemon
Documentation=https://atuin.sh/docs
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Use correct path - check both locations
ExecStartPre=/bin/sh -c 'test -x /usr/bin/atuin'
ExecStart=usr/bin/atuin daemon
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
Environment=ATUIN_CONFIG_DIR=%h/.config/atuin
KillMode=mixed
TimeoutStopSec=5

# Security settings (relaxed for proper operation)
ProtectSystem=false
PrivateTmp=true
PrivateDevices=false
ProtectHome=false
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictRealtime=true
LockPersonality=true
RestrictSUIDSGID=true
NoNewPrivileges=true
ProtectHostname=true
ProtectClock=true

[Install]
WantedBy=default.target
