# .config/systemd/user/hyprtheme.service
[Unit]
Description=Apply GTK theme configuration for Hyprland
After=graphical-session.target dbus.service
Wants=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/usr/bin/bash -c 'systemctl --user import-environment DISPLAY WAYLAND_DISPLAY XDG_CURRENT_DESKTOP'
ExecStart=/usr/bin/bash -c 'set -euo pipefail; \
    sleep 5; \
    export THEME="Colloid-Dark-Compact-Catppuccin"; \
    export ICONS="Papirus"; \
    export FONT="JetBrainsMono Nerd Font Medium 13"; \
    export CURSOR="Adwaita"; \
    if [[ -z "${DBUS_SESSION_BUS_ADDRESS:-}" ]]; then \
        export DBUS_SESSION_BUS_ADDRESS="unix:path=${XDG_RUNTIME_DIR}/bus"; \
    fi; \
    if ! command -v gsettings &>/dev/null; then \
        logger -t hyprtheme "Error: gsettings command not found"; \
        exit 1; \
    fi; \
    timeout 10 gsettings set org.gnome.desktop.interface gtk-theme "$THEME" || { logger -t hyprtheme "Failed to set gtk-theme"; exit 1; }; \
    timeout 10 gsettings set org.gnome.desktop.interface icon-theme "$ICONS" || { logger -t hyprtheme "Failed to set icon-theme"; exit 1; }; \
    timeout 10 gsettings set org.gnome.desktop.interface cursor-theme "$CURSOR" || { logger -t hyprtheme "Failed to set cursor-theme"; exit 1; }; \
    timeout 10 gsettings set org.gnome.desktop.interface font-name "$FONT" || { logger -t hyprtheme "Failed to set font-name"; exit 1; }; \
    logger -t hyprtheme "Applied theme configuration successfully"'

Restart=on-failure
RestartSec=15
TimeoutStartSec=90

# Import user session environment
PassEnvironment=DISPLAY WAYLAND_DISPLAY XDG_CURRENT_DESKTOP XDG_SESSION_TYPE XDG_SESSION_DESKTOP

# Standard environment variables
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=XDG_DATA_DIRS=/usr/local/share:/usr/share:/var/lib/flatpak/exports/share:/home/%i/.local/share/flatpak/exports/share
Environment=XDG_CONFIG_DIRS=/etc/xdg
Environment=XDG_CACHE_HOME=/home/%i/.cache
Environment=XDG_CONFIG_HOME=/home/%i/.config
Environment=XDG_DATA_HOME=/home/%i/.local/share
Environment=XDG_STATE_HOME=/home/%i/.local/state
Environment=XDG_RUNTIME_DIR=/run/user/%i

# Security hardening (balanced for desktop integration)
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/%i/.config /home/%i/.local/share /home/%i/.cache /run/user/%i
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RemoveIPC=yes
ProtectHostname=yes
ProtectClock=yes
ProtectKernelLogs=yes
SystemCallArchitectures=native
LockPersonality=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
SystemCallFilter=@system-service @signal @ipc @sync @file-system @process @network-io @basic-io
SystemCallErrorNumber=EPERM

# Capabilities (none needed)
CapabilityBoundingSet=
AmbientCapabilities=

# Resource limits
TasksMax=50
MemoryMax=200M
CPUQuota=30%

[Install]
WantedBy=graphical-session.target
