# .config/systemd/user/tldr-update.service
[Unit]
Description=Update tldr cache
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/tldr --update
